# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-176ccbe8\x2d961e\x2d4b77\x2dbbd0\x2d40bee15e7fa5.target

[Mount]
Where=/
What=/dev/disk/by-uuid/176ccbe8-961e-4b77-bbd0-40bee15e7fa5
Type=ext4
