# etc/udev/rules.d/99-vultr-fix-virtio.rules
KERNEL=="vd*[!0-9]" \
, ACTION=="add" \
, RUN+="/usr/sbin/udevadm trigger --name-match=$env{DEVNAME}"
KERNEL=="vd*[0-9]" \
, ACTION=="add" \
, RUN+="/usr/sbin/udevadm trigger --name-match=$env{DEVNAME}"